# crontab
*/1 * * * * /opt/python-rs1168/stats-sr1168.py 
