# crontab
*/1 * * * * /opt/sr1168c/stats-sr1168c.py 
