# crontab
*/1 * * * * /opt/python-rs1168/influx-sr1168.py 
